# bin/var/lib/dpkg/info/mercurial-common.conffiles
/etc/mercurial/hgrc
/etc/mercurial/hgrc.d/cacerts.rc
